# scr-protector/systemd/scr-protector-blocker.timer
[Unit]
Description=scr-protector Blocklist Sync Timer
Documentation=https://github.com/yourusername/scr-protector
Requires=scr-protector-blocker.service

[Timer]
# Start 10 seconds after boot
OnBootSec=10s

# Run every 30 seconds
OnUnitActiveSec=30s

# Add randomized delay to prevent thundering herd
RandomizedDelaySec=5s

# Persist timer across reboots
Persistent=true

[Install]
WantedBy=timers.target
