# scr-protector/systemd/scr-protector-parser.service
[Unit]
Description=scr-protector Log Parser and Threat Detector
Documentation=https://github.com/yourusername/scr-protector
After=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/scr-protector
Environment="PATH=/opt/scr-protector/venv/bin"

# Run the parser daemon
ExecStart=/opt/scr-protector/venv/bin/python /opt/scr-protector/core/parser.py

# Restart on failure
Restart=on-failure
RestartSec=5

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/scr-protector
ReadOnlyPaths=/var/log/nginx /var/log/auth.log

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=scr-parser

[Install]
WantedBy=multi-user.target
