# scr-protector/systemd/scr-protector-blocker.service
[Unit]
Description=scr-protector Blocklist Sync
Documentation=https://github.com/yourusername/scr-protector
After=network.target

[Service]
Type=oneshot
User=root
Group=root

# Run the blocklist sync script
ExecStart=/opt/scr-protector/bin/apply_blocklist.sh

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=scr-blocker
